# services/systemd/user/stock-db-downstream-refresh.service
[Unit]
Description=Refresh downstream project JSON data after price update
After=stock-db-price-refresh.service

[Service]
Type=oneshot
WorkingDirectory=/home/exp/projects/stock_db
Environment=PATH=/etc/profiles/per-user/exp/bin:/run/current-system/sw/bin:/usr/bin:/bin
ExecStart=/etc/profiles/per-user/exp/bin/uv run python services/downstream_refresh.py
